# 99-nintendo-licensed-led.rules
ACTION=="add", SUBSYSTEM=="hidraw", KERNELS=="0005:0E6F:018C.*", RUN+="/usr/bin/procon-led /dev/%k"
